# deployments/native/linux/sonarqube.service.in
[Unit]
Description=SonarQube installed by SonarWeaver
Documentation=https://github.com/Yunushan/SonarWeaver
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=@SERVICE_USER@
Group=@SERVICE_GROUP@
ExecStart=/usr/local/libexec/sonarweaver-start
Restart=on-failure
RestartSec=10
TimeoutStartSec=300
TimeoutStopSec=120
LimitNOFILE=131072
LimitNPROC=8192
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectKernelLogs=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true
UMask=0027

[Install]
WantedBy=multi-user.target
